# code-backup.service
[Unit]
Description=Code Backup Daemon with Web UI
Documentation=https://github.com/your-username/code-backup-daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/nayan-ai4m/Desktop/NK/automated-github-push
ExecStart=/home/nayan-ai4m/Desktop/NK/automated-github-push/venv/bin/python -m code_backup_daemon.cli start
Restart=on-failure
RestartSec=30
TimeoutStopSec=30

# Environment variables
Environment=HOME=/home/nayan-ai4m
EnvironmentFile=/home/nayan-ai4m/Desktop/NK/automated-github-push/.env

# Security settings (relaxed for user service)
ProtectHome=false

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=code-backup-daemon

[Install]
WantedBy=default.target
