# code-backup.service
[Unit]
Description=Code Backup Daemon
Documentation=https://github.com/your-username/code-backup-daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=%h
ExecStart=/usr/local/bin/code-backup start
ExecStop=/usr/local/bin/code-backup stop
Restart=on-failure
RestartSec=30
TimeoutStopSec=30

# Environment variables
Environment=HOME=%h
Environment=USER=%i

# Security settings
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=%h/.local/share/code-backup %h/.config/code-backup %h/CODE

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=code-backup-daemon

[Install]
WantedBy=default.target
